# com.noslip.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noslip.telegram</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/your/cloned/repository/services/trader/.venv/bin/python</string>
        <string>-u</string>
        <string>/path/to/your/cloned/repository/services/trader/telegram_interactive_bot.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/path/to/your/cloned/repository</string>
    <key>StandardOutPath</key>
    <string>/path/to/your/cloned/repository/data/telegram_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/your/cloned/repository/data/telegram_bot.log</string>
</dict>
</plist>
